# tracking-bot.service
[Unit]
Description=Tracking Bot
After=network.target

[Service]
Type=simple
User=your_user
WorkingDirectory=/home/your_user/tracking-bot
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/home/your_user/tracking-bot

[Install]
WantedBy=multi-user.target
